# roles/logger-app-3/templates/logger_app_3.service.j2
[Unit]
Description=Logger APP (3)

[Service]
Type=simple
ExecStart=/bin/bash /opt/logger_app_3/logger_app_3.sh
Restart=always
User=nobody

[Install]
WantedBy=multi-user.target

